# ansible/templates/systemd/minecraft.service.j2
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=forking
{% if random_server_icon %}
ExecStartPre=-/opt/minecraft/server-icon.sh
{% endif %}
ExecStart=/usr/bin/tmux new-session -s minecraft -d '{{ command }}'
ExecStop=/usr/bin/tmux send-keys -t minecraft stop ENTER

WorkingDirectory=/opt/minecraft/server

User=minecraft
Group=minecraft

Nice=-20
OOMScoreAdjust=-1000

[Install]
WantedBy=multi-user.target
