# ansible/templates/systemd/minecraft-backup.timer.j2
[Unit]
Description=Backup Minecraft {{ edition | capitalize }} Server

[Timer]
OnCalendar=6:00
Persistent=true

[Install]
WantedBy=timers.target
