# ansible/templates/systemd/minecraft.service.j2
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=forking
{% if ansible_memtotal_mb < (12288 + heap) %}
ExecStart=/usr/bin/tmux new-session -s minecraft -d '{{ command }}'
{% else %}
ExecStart=/usr/bin/tmux new-session -s minecraft -d '{{ high_ram_command }}'
{% endif %}
ExecStop=/usr/bin/tmux send-keys -t minecraft C-c

WorkingDirectory=/opt/minecraft/server

User=minecraft
Group=minecraft

[Install]
WantedBy=multi-user.target
