# ansible/templates/systemd/minecraft.service.j2
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=forking
{% if ansible_memtotal_mb < (12288 + 1024) %}
ExecStart=/usr/bin/tmux new-session -s minecraft -d '{{ command }}'
{% else %}
ExecStart=/usr/bin/tmux new-session -s minecraft -d '{{ high_ram_command }}'
{% endif %}
ExecStop=/opt/minecraft/stop_minecraft.sh

WorkingDirectory=/opt/minecraft/server

User=minecraft
Group=minecraft

[Install]
WantedBy=multi-user.target
