# ansible/templates/systemd/minecraft.service.j2
[Unit]
Description=Minecraft {{ edition | capitalize }} Server
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/tmux new-session -s minecraft-{{ edition }} -d '{{ command }}'
{% if type != "bedrockconnect" -%}
ExecStop=/home/minecraft/stop_{{ edition }}.sh
{% endif %}

WorkingDirectory=/home/minecraft/{{ edition }}

User=minecraft
Group=minecraft

[Install]
WantedBy=multi-user.target
