# ansible/templates/systemd/monitor.timer.j2
[Unit]
Description={{ item.description }} Monitor

[Timer]
OnCalendar={{ item.timer }}

[Install]
WantedBy=timers.target
