# ansible/files/systemd/apt-upgrade.timer
[Unit]
Description=Apt Upgrade

[Timer]
OnCalendar=Tue 7:00
Persistent=true

[Install]
WantedBy=timers.target
