# ansible/files/systemd/heap-utilization-monitor.timer
[Unit]
Description=Heap Usage Monitor

[Timer]
OnCalendar=*-*-* *:*:*
AccuracySec=1s

[Install]
WantedBy=timers.target
