# ansible/files/systemd/minecraft-bedrock-backup.timer
[Unit]
Description=Backup Minecraft server

[Timer]
OnCalendar=1:00
Persistent=true

[Install]
WantedBy=timers.target
