# ansible/files/systemd/oci-disk-utilization.timer
[Unit]
Description=OCI Disk Utilization Monitor

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
