# ansible/files/systemd/oci-disk-utilization.timer
[Unit]
Description=OCI Disk Utilization Monitor

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s

[Install]
WantedBy=timers.target
