# ansible/files/systemd/query-monitor.timer
[Unit]
Description=Minecraft Server Status

[Timer]
OnCalendar=*-*-* *:*:10
AccuracySec=1m

[Install]
WantedBy=timers.target
