# ansible/files/systemd/query-monitor.timer
[Unit]
Description=Minecraft Server Status

[Timer]
OnCalendar=*:0/10:0
AccuracySec=1m

[Install]
WantedBy=timers.target
