# ansible/files/systemd/tick-monitor.timer
[Unit]
Description=Minecraft MSPT

[Timer]
OnCalendar=*:*:*
AccuracySec=1m

[Install]
WantedBy=timers.target
